# roles/wireguard/tasks/files/wireguard_reresolve-dns.timer
[Unit]
Description=Periodically reresolve DNS of all WireGuard endpoints

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
