# roles/proxmox-backup/templates/proxmox-autosnap-weekly.timer
[Unit]
Description=Run proxmox-autosnap weekly

[Timer]
OnCalendar=Mon *-*-2..31 00:05:00

[Install]
WantedBy=timers.target
